# backup/crontab
0 2 * * * /backup.sh >> /var/log/backup.log 2>&1
